# build/configs/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log

[program:adminpanel]
command=/var/www/onlyoffice/documentserver/server/AdminPanel/server/adminpanel
directory=/var/www/onlyoffice/documentserver/server/AdminPanel
user=ds
environment=NODE_ENV=development-linux,NODE_CONFIG_DIR=/etc/onlyoffice/documentserver,NODE_DISABLE_COLORS=1,APPLICATION_NAME=onlyoffice,LD_LIBRARY_PATH=/var/www/onlyoffice/documentserver/server/FileConverter/bin:/var/www/onlyoffice/documentserver/server/FileConverter/lib
stdout_logfile=/var/log/onlyoffice/documentserver/adminpanel/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/adminpanel/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:converter]
command=/var/www/onlyoffice/documentserver/server/FileConverter/converter
directory=/var/www/onlyoffice/documentserver/server/FileConverter
user=ds
environment=NODE_ENV=development-linux,NODE_CONFIG_DIR=/etc/onlyoffice/documentserver,NODE_DISABLE_COLORS=1,APPLICATION_NAME=onlyoffice,LD_LIBRARY_PATH=/var/www/onlyoffice/documentserver/server/FileConverter/bin:/var/www/onlyoffice/documentserver/server/FileConverter/lib
stdout_logfile=/var/log/onlyoffice/documentserver/converter/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/converter/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:docservice]
command=/var/www/onlyoffice/documentserver/server/DocService/docservice
directory=/var/www/onlyoffice/documentserver/server/DocService
user=ds
environment=NODE_ENV=development-linux,NODE_CONFIG_DIR=/etc/onlyoffice/documentserver,NODE_DISABLE_COLORS=1,APPLICATION_NAME=onlyoffice,LD_LIBRARY_PATH=/var/www/onlyoffice/documentserver/server/FileConverter/bin:/var/www/onlyoffice/documentserver/server/FileConverter/lib
stdout_logfile=/var/log/onlyoffice/documentserver/docservice/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/docservice/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:metrics]
command=/var/www/onlyoffice/documentserver/server/Metrics/metrics ./config/config.js
directory=/var/www/onlyoffice/documentserver/server/Metrics
user=ds
environment=NODE_DISABLE_COLORS=1,LD_LIBRARY_PATH=/var/www/onlyoffice/documentserver/server/FileConverter/bin:/var/www/onlyoffice/documentserver/server/FileConverter/lib
stdout_logfile=/var/log/onlyoffice/documentserver/metrics/out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/metrics/err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:ds-example]
command=/var/www/onlyoffice/documentserver-example/example
directory=/var/www/onlyoffice/documentserver-example/
user=ds
environment=NODE_ENV=development-linux,NODE_CONFIG_DIR=/etc/onlyoffice/documentserver-example,NODE_DISABLE_COLORS=1,APPLICATION_NAME=onlyoffice
stdout_logfile=/var/log/onlyoffice/documentserver/ds-example_out.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/onlyoffice/documentserver/ds-example_err.log
stderr_logfile_backups=0
stderr_logfile_maxbytes=0
autostart=false
